# client/config/kiosk-launcher.service
# systemd unit для Kiosk Launcher
# Запуск при загрузке системы, автоперезапуск при падении.
#
# Установка (от root или через sudo):
#   cp kiosk-launcher.service /etc/systemd/system/
#   # Отредактируйте Environment и пути в [Service] при необходимости
#   systemctl daemon-reload
#   systemctl enable kiosk-launcher
#   systemctl start kiosk-launcher
#
# Управление:
#   systemctl status kiosk-launcher
#   systemctl stop kiosk-launcher
#   journalctl -u kiosk-launcher -f

[Unit]
Description=Kiosk Launcher (Wi-Fi + app watchdog)
After=network.target NetworkManager.service
Wants=network-online.target

[Service]
Type=simple
# Путь к скрипту и рабочая директория (где лежит приложение)
WorkingDirectory=/home/live/kiosk-app
ExecStart=/home/live/kiosk-app/kiosk-launcher.sh
# Автоперезапуск при падении (в т.ч. при ошибке в скрипте)
Restart=always
RestartSec=10
# От какого пользователя запускать (киоск обычно один пользователь)
User=live
Group=live

[Install]
WantedBy=multi-user.target
